# ttx_hbase/templates/default/hbase-initd.erb
#!/bin/bash
# chkconfig: 2345 95 20
# description: hbase-<%= @script_service %>
# processname: java


#
# Copyright 2013 ConnectDotz, LLC. All Rights Reserved.
#

. /etc/rc.d/init.d/functions
. <%= node[:ttx_hbase][:hbase][:conf] %>/hbase-env.sh

pid=$HBASE_PID_DIR/hbase-$HBASE_IDENT_STRING-<%= @script_service %>.pid

cd $HBASE_HOME

# set the ulimit to 32768 because hbase needs lots of open files
ulimit -n 32768

function start {
    su -l -c ". <%= node[:ttx_hbase][:hbase][:conf] %>/hbase-env.sh; $HBASE_HOME/bin/hbase-daemon.sh --config $HBASE_CONF_DIR start <%= @script_service %>  <%= @script_args %>"  
}

function stop {
    su -l -c ". <%= node[:ttx_hbase][:hbase][:conf] %>/hbase-env.sh; $HBASE_HOME/bin/hbase-daemon.sh --config $HBASE_CONF_DIR stop <%= @script_service %> <%= @script_args %>" 
}

function reload {
    su -l -c ". <%= node[:ttx_hbase][:hbase][:conf] %>/hbase-env.sh; $HBASE_HOME/bin/hbase-daemon.sh --config $HBASE_CONF_DIR restart <%= @script_service %> <%= @script_args %>" 
}

function status {
#    log_warning_msg "Status is not supported"
	if [ -f $pid ]; then
	    echo "<%= @script_service %> is already running"
		exit 0
	else
		echo "<%= @script_service %> is not running"
		exit 1
	fi
}

case $1 in
    'start' )
	start
	;;
    'stop' )
	stop
	;;
    'restart' )
	stop
	start
	;;
    'force-reload' )
	reload
	;;
    'status' )
	status
	;;
    *)
	echo "usage: `basename $0` {start|stop|status}"
esac 

exit 0 

